# src/server/minecraft/minecraft@.socket
[Unit]
PartOf=minecraft@%i.service

[Socket]
ListenFIFO=%t/minecraft.stdin
