# src/server/systemd/timedrestart.service
[Unit]
Description=Timed restart script

[Service]
Type=oneshot
ExecStart=echo "Nightly naptime! Rebooting..."
ExecStart=/usr/bin/reboot
